# systemd/user/qwen3-npu-aichat.service
[Unit]
Description=Local Qwen3 NPU OpenAI bridge
After=default.target

[Service]
Type=simple
ExecStart=/home/lukas/drivers/qwen3-8b-fresh/run-openai-server.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
